# container-files/etc/supervisor.d/bitbucket.conf
[program:bitbucket]
command = /opt/atlassian/bitbucket/start-bitbucket.sh -fg
autorestart = true
