# container-files/etc/supervisor.d/bitbucket.conf
[program:bitbucket]
command = /opt/atlassian/bitbucket/bin/start-jira.sh -fg
autorestart = true
